# supervisord.conf
[supervisord]
nodaemon=true

[program:tinc]
command=/root/scripts/run.sh
autostart=true
autorestart=true
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
